# logtestunit.timer
[Unit]
Description=timer for logging

[Timer]
OnCalendar=*:*:00

[Install]
WantedBy=timers.target