# logtestunit.service
[Unit]
Description=log script date who uptime

[Service]
ExecStart=/usr/MFI/logging